# etc/metricshipper_supervisor.conf
[program:metricshipper]
command=/opt/zenoss/bin/metricshipper -c /opt/zenoss/etc/metricshipper/metricshipper.yaml
autorestart=true
autostart=true
startsecs=5
stopwaitsecs=20

; logging
redirect_stderr=true
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stdout_logfile=/opt/zenoss/log/%(program_name)s.log
